# debian/postinst
#!/bin/sh
set -e

# Detect config.txt location
if [ -f "/boot/firmware/config.txt" ]; then
    CONFIG_FILE="/boot/firmware/config.txt"
elif [ -f "/boot/config.txt" ]; then
    CONFIG_FILE="/boot/config.txt"
else
    echo "Warning: config.txt not found in /boot/firmware or /boot, skipping dtoverlay configuration" >&2
    CONFIG_FILE=""
fi

SERVICE_FILE="/etc/systemd/system/usb-ethernet-gadget.service"

# Make usb-ethernet-gadget.sh executable
chmod +x /usr/bin/usb-ethernet-gadget.sh

if [ -n "$CONFIG_FILE" ]; then
    # Comment out any existing otg_mode=1 entries
    sed -i 's/^otg_mode=1/# otg_mode=1/' "$CONFIG_FILE"

    # Idempotent config.txt processing:
    # 1. Check if [cm4] section exists
    # 2. Check if dtoverlay=dwc2 already exists in [cm4] section
    # 3. Only add if missing

    if grep -q '^\[cm4\]' "$CONFIG_FILE"; then
        # [cm4] section exists - check if dtoverlay=dwc2 is already there
        # Extract [cm4] section and check for dtoverlay=dwc2
        if ! awk '/^\[cm4\]/,/^\[/' "$CONFIG_FILE" | grep -q '^dtoverlay=dwc2'; then
            # dtoverlay=dwc2 not in [cm4] section - add it after [cm4] header
            sed -i '/^\[cm4\]/a dtoverlay=dwc2' "$CONFIG_FILE"
            echo "Added dtoverlay=dwc2 to existing [cm4] section"
        else
            echo "dtoverlay=dwc2 already present in [cm4] section"
        fi
    else
        # No [cm4] section - add it at the end of the file
        printf '\n[cm4]\ndtoverlay=dwc2\n' >> "$CONFIG_FILE"
        echo "Added [cm4] section with dtoverlay=dwc2"
    fi
fi

# Reload systemd configuration
systemctl daemon-reload

# Stop and disable old usb-gadget.service if it exists
if [ -f "/etc/systemd/system/usb-gadget.service" ]; then
    systemctl stop usb-gadget.service || true
    systemctl disable usb-gadget.service || true
fi

# Update /etc/avahi/avahi-daemon.conf for USB Ethernet
AVAHI_CONF="/etc/avahi/avahi-daemon.conf"
if [ -f "$AVAHI_CONF" ]; then
    echo "Updating $AVAHI_CONF for USB Ethernet"

    # Insert or update "allow-interfaces=usb0,usb1,wlan0,pan0"
    if grep -q '^allow-interfaces=' "$AVAHI_CONF"; then
        sed -i 's/^allow-interfaces=.*/allow-interfaces=usb0,usb1,wlan0,pan0/' "$AVAHI_CONF"
    else
        sed -i '/^\[server\]/a allow-interfaces=usb0,usb1,wlan0,pan0' "$AVAHI_CONF"
    fi

    # Insert or update "publish-hinfo=no"
    if grep -q '^publish-hinfo=' "$AVAHI_CONF"; then
        sed -i 's/^publish-hinfo=.*/publish-hinfo=no/' "$AVAHI_CONF"
    else
        sed -i '/^\[server\]/a publish-hinfo=no' "$AVAHI_CONF"
    fi
else
    echo "Warning: $AVAHI_CONF not found; unable to configure Avahi for USB Ethernet."
fi

# Detect if running in a chroot environment
if systemd-detect-virt --chroot >/dev/null 2>&1; then
    echo "Running in a chroot environment. Only enabling services."
    # Enable services without starting them
    systemctl enable systemd-networkd || true
    if [ -f "$SERVICE_FILE" ]; then
        systemctl enable usb-ethernet-gadget.service || true
    else
        echo "Error: usb-ethernet-gadget.service not found at $SERVICE_FILE" >&2
        exit 1
    fi
else
    # Enable and start services in a normal environment
    systemctl enable systemd-networkd
    systemctl start systemd-networkd || true
    if [ -f "$SERVICE_FILE" ]; then
        systemctl enable usb-ethernet-gadget.service
        systemctl start usb-ethernet-gadget.service || true
    else
        echo "Error: usb-ethernet-gadget.service not found at $SERVICE_FILE" >&2
        exit 1
    fi

    # Optionally restart Avahi so changes take effect immediately
    if systemctl is-enabled avahi-daemon >/dev/null 2>&1; then
        systemctl restart avahi-daemon || true
    fi
fi

exit 0
