# etc/systemd/system/usb-ethernet-gadget.service
[Unit]
Description=USB ethernet gadget setup and keep-alive service
After=local-fs.target systemd-modules-load.service sys-kernel-config.mount
Before=network-pre.target
DefaultDependencies=no
RequiresMountsFor=/var /sys/kernel/config
ConditionPathExists=/usr/bin/usb-ethernet-gadget.sh

[Service]
Type=simple
ExecStart=/usr/bin/usb-ethernet-gadget.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
