# etc/systemd/system/usb-ethernet-gadget.service
[Unit]
Description=USB Ethernet Gadget Setup and Keep-Alive Service
After=local-fs.target systemd-modules-load.service
Before=network-pre.target
DefaultDependencies=no
RequiresMountsFor=/var

[Service]
Type=simple
ExecStart=/usr/bin/usb-ethernet-gadget.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
